# ansible/phases/05-5g-core/roles/runtime_watchdog/templates/k3s-pull-watchdog.service.j2
[Unit]
Description=K3s image pull watchdog (auto-restart on repeated failures)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart={{ watchdog_script_path }}
Nice=5
Environment=THRESHOLD={{ watchdog_threshold }}


